# pishoo/pkg/debian/postinst
#!/bin/sh
set -e

LOG_DIR="/var/log/pishoo"
LOG_PERMISSIONS="755"

case "$1" in
  configure)
        if [ ! -d "$LOG_DIR" ]; then
            mkdir -p "$LOG_DIR"
            if [ $? -ne 0 ]; then
                echo "Failed to create log directory $LOG_DIR" >&2
                exit 1
            fi

            chmod $LOG_PERMISSIONS "$LOG_DIR"
            if [ $? -ne 0 ]; then
                echo "Failed to set permissions on log directory $LOG_DIR" >&2
                exit 1
            fi
          fi
        fi
    ;;
  triggered)
    if command -v systemctl > /dev/null 2>&1 && systemctl is-active --quiet pishoo; then
      echo "Triggering pishoo reload ..."
      systemctl reload pishoo
    elif [ -x "/etc/init.d/pishoo" ] && invoke-rc.d --quiet pishoo status >/dev/null; then
      echo "Triggering pishoo reload ..."
      invoke-rc.d pishoo reload || true
    fi
    exit 0
    ;;
  *)
    echo "postinst called with unknown argument \`$1'" >&2
    ;;
esac

#DEBHELPER#

if command -v systemctl > /dev/null 2>&1 && systemctl is-active --quiet pishoo; then
  systemctl upgrade pishoo || systemctl restart pishoo
  exit $?
elif [ -x "/etc/init.d/pishoo" ] && invoke-rc.d --quiet pishoo status >/dev/null; then
  invoke-rc.d pishoo upgrade || invoke-rc.d pishoo restart
  exit $?
else
  echo "Not attempting to start Pishoo, please start Pishoo yourself."
  exit 0
fi


exit 0